# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/07_mixed_deps_groups/b.service
[Unit]
Description=b.service UNIT FILE
Before=a.service g1.target

[Service]
ExecStart=/usr/bin/b

[Install]
RequiredBy=a.service g1.target

