# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/07_mixed_deps_groups/g1.target
[Unit]
Description=g1.target UNIT FILE
Requires=a.service b.service
After=a.service b.service

