# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/08_group_in_group/g1.target
[Unit]
Description=g1.target UNIT FILE
Before=g2.target

[Install]
RequiredBy=g2.target

