# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/08_group_in_group/g2.target
[Unit]
Description=g2.target UNIT FILE
Requires=g1.target
After=g1.target

