# testing/UnitFileGeneration/PODMAN_COMPOSE/TEST_CASES/08_abort_cascade/svc.service
[Unit]
Description=svc.service UNIT FILE
Requires=p.service
PartOf=p.service
After=p.service

[Service]
Type=simple
ExecStart=podman run --name svc docker.io/library/alpine:latest
ExecStop=podman stop svc
ExecStopPost=podman rm svc

