# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/10_abort_ignore/c.service
[Unit]
Description=c.service UNIT FILE
Requires=p.service
After=p.service

[Service]
ExecStart=/usr/bin/c

