# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/06_nested_groups/a.service
[Unit]
Description=a.service UNIT FILE
Before=g1.target

[Service]
ExecStart=/usr/bin/a

[Install]
RequiredBy=g1.target

