# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/06_nested_groups/b.service
[Unit]
Description=b.service UNIT FILE
Before=g1.target

[Service]
ExecStart=/usr/bin/b

[Install]
RequiredBy=g1.target

