# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/06_nested_groups/g1.target
[Unit]
Description=g1.target UNIT FILE
Requires=a.service b.service
Before=g2.target
After=a.service b.service

[Install]
RequiredBy=g2.target

