# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/03_entry_with_deps/a.service
[Unit]
Description=a.service UNIT FILE
Requires=b.service
After=b.service

[Service]
ExecStart=/usr/bin/a

