# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/03_entry_with_deps/b.service
[Unit]
Description=b.service UNIT FILE
Requires=c.service
Before=a.service
After=c.service

[Service]
ExecStart=/usr/bin/b

[Install]
RequiredBy=a.service

