# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/03_entry_with_deps/c.service
[Unit]
Description=c.service UNIT FILE
Before=b.service

[Service]
ExecStart=/usr/bin/c

[Install]
RequiredBy=b.service

