# testing/UnitFileGeneration/PODMAN_COMPOSE/TEST_CASES/11_cascade_multi_parent/p1.service
[Unit]
Description=p1.service UNIT FILE

[Service]
Type=simple
ExecStart=podman run --name p1 docker.io/library/alpine:latest
ExecStop=podman stop p1
ExecStopPost=podman rm p1

