# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/14_cascade_multi_parent/c.service
[Unit]
Description=c.service UNIT FILE
Requires=p1.service p2.service
PartOf=p1.service p2.service
After=p1.service p2.service
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
ExecStart=/usr/bin/c
RestartSec=2
Restart=on-failure

