# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/14_cascade_multi_parent/p1.service
[Unit]
Description=p1.service UNIT FILE

[Service]
ExecStart=/usr/bin/p1

