# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/14_cascade_multi_parent/p2.service
[Unit]
Description=p2.service UNIT FILE

[Service]
ExecStart=/usr/bin/p2

