# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/01_minimal_entry/a.service
[Unit]
Description=a.service UNIT FILE

[Service]
ExecStart=/usr/bin/a

