# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/09_large_mixed/a.service
[Unit]
Description=a.service UNIT FILE
Requires=b.service c.service
Before=g1.target
After=b.service c.service

[Service]
ExecStart=/usr/bin/a

[Install]
RequiredBy=g1.target

