# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/09_large_mixed/e.service
[Unit]
Description=e.service UNIT FILE
Before=d.service f.service g2.target

[Service]
ExecStart=/usr/bin/e

[Install]
RequiredBy=d.service f.service g2.target

