# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/09_large_mixed/g2.target
[Unit]
Description=g2.target UNIT FILE
Requires=d.service e.service f.service g1.target
After=d.service e.service f.service g1.target

