# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/09_large_mixed/i.service
[Unit]
Description=i.service UNIT FILE

[Service]
ExecStart=/usr/bin/i /tmp/result.log

