# testing/UnitFileGeneration/PODMAN_COMPOSE/TEST_CASES/04_service_in_network/appnet.service
[Unit]
Description=appnet.service UNIT FILE
Requires=web.service
After=web.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=podman network create --ignore appnet
ExecStop=podman network rm appnet

