# testing/UnitFileGeneration/PODMAN_COMPOSE/TEST_CASES/05_network_with_members_and_depends/appnet.service
[Unit]
Description=appnet.service UNIT FILE
Requires=basenet.service web.service db.service
After=basenet.service web.service db.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=podman network create --ignore appnet
ExecStop=podman network rm appnet

