# testing/UnitFileGeneration/PODMAN_COMPOSE/TEST_CASES/05_network_with_members_and_depends/basenet.service
[Unit]
Description=basenet.service UNIT FILE
Before=appnet.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=podman network create --ignore basenet
ExecStop=podman network rm basenet

[Install]
RequiredBy=appnet.service

