# testing/UnitFileGeneration/PODMAN_COMPOSE/TEST_CASES/05_network_with_members_and_depends/db.service
[Unit]
Description=db.service UNIT FILE
Before=appnet.service

[Service]
Type=simple
ExecStart=podman run --name db docker.io/library/alpine:latest
ExecStop=podman stop db
ExecStopPost=podman rm db

[Install]
RequiredBy=appnet.service

