# testing/UnitFileGeneration/PODMAN_COMPOSE/TEST_CASES/05_network_with_members_and_depends/web.service
[Unit]
Description=web.service UNIT FILE
Before=appnet.service

[Service]
Type=simple
ExecStart=podman run --name web docker.io/library/alpine:latest
ExecStop=podman stop web
ExecStopPost=podman rm web

[Install]
RequiredBy=appnet.service

