# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/02_minimal_group/g1.target
[Unit]
Description=g1.target UNIT FILE

