# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/13_restart_cascade/p.service
[Unit]
Description=p.service UNIT FILE

[Service]
ExecStart=/usr/bin/p

