# testing/UnitFileGeneration/GENERIC_TASK_LAUNCHER/TEST_CASES/13_restart_cascade/q.service
[Unit]
Description=q.service UNIT FILE

[Service]
ExecStart=/usr/bin/q

